# deploy/polymomentum-engine.service
[Unit]
Description=PolyMomentum Rust trading engine (paper/live)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polymomentum
WorkingDirectory=/opt/polymomentum
EnvironmentFile=/etc/polymomentum/env
Environment=RUST_LOG=info
Environment=POLYMOMENTUM_DATA_DIR=/opt/polymomentum/data
Environment=POLYMOMENTUM_LOGS_DIR=/opt/polymomentum/logs
Environment=STATE_DB_PATH=/opt/polymomentum/logs/candle/state.db
Environment=SESSION_LOG_DIR=/opt/polymomentum/logs/sessions
Environment=KILL_SWITCH_PATH=/tmp/polymomentum/KILL
Environment=BANKROLL_USD=100
ExecStartPre=/opt/polymomentum/polymomentum-engine preflight --mode paper
ExecStart=/opt/polymomentum/polymomentum-engine live --mode paper
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Nice=5
CPUQuota=80%
MemoryMax=512M
TasksMax=256
IOSchedulingClass=best-effort
IOSchedulingPriority=6

# Hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/polymomentum /tmp/polymomentum
ReadOnlyPaths=-/opt/shared

[Install]
WantedBy=multi-user.target
